# bacula/supervisord.conf
[supervisord]
nodaemon=true


[program:bacula]
command=/usr/local/bacula/sbin/bacula-fd -f



