# haproxy/supervisord.conf
[supervisord]
nodaemon=true


[program:haproxy]
command=/usr/local/sbin/haproxy -f /etc/haproxy/haproxy.cfg
